# launchd/io.rook.agent.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.rook.agent</string>

  <key>ProgramArguments</key>
  <array>
    <string>__REPO_ROOT__/scripts/run-rook.sh</string>
    <string>__CONFIG_PATH__</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>__LOG_DIR__/rook.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/rook.stderr.log</string>
</dict>
</plist>
